# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
childlogdir=/var/log/supervisor

[program:redis]
# Save if ≥1 key changes within 30 s (matches the upstream docker-compose setting).
# Adjust the --save interval as needed, or pass --save "" to disable persistence.
command=redis-server --save 30 1 --loglevel warning --dir /data
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log

[program:backend]
command=python3 /app/api.py
directory=/app
user=root
autostart=true
autorestart=true
# No environment= directive here so the process inherits the full container
# environment, including REDIS_URL (default set in Dockerfile) as well as all
# user-supplied variables: PROVIDER_NAME, PROVIDER_MODEL,
# PROVIDER_SERVER_ADDRESS, SEARXNG_BASE_URL, API keys, etc.
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-err.log

[program:frontend]
command=npm start
directory=/frontend
user=root
autostart=true
autorestart=true
# No environment= directive here so the process inherits the full container
# environment, including PORT and CHOKIDAR_USEPOLLING (defaults set in
# Dockerfile) as well as REACT_APP_BACKEND_URL supplied by the user.
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend-err.log
